# scripts/nightly/com.edihasaj.vmlab.nightly.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.edihasaj.vmlab.nightly</string>

    <key>ProgramArguments</key>
    <array>
        <string>__VMLAB_BIN__</string>
        <string>run</string>
        <string>@@vm</string>
        <string>examples/flows/recall-cross-os.yaml</string>
        <string>--max-parallel</string>
        <string>2</string>
    </array>

    <key>WorkingDirectory</key>
    <string>__REPO__</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>__LOG_DIR__/nightly.out.log</string>
    <key>StandardErrorPath</key>
    <string>__LOG_DIR__/nightly.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>__PATH__</string>
        <key>HOME</key>
        <string>__HOME__</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
